# ticket-runner-dashboard.service
[Unit]
Description=Ticket runner dashboard (read-only view of state/runner.db)
Wants=network-online.target
After=network-online.target ticket-runner.service
# Stops/restarts together with the runner via the shared target.
PartOf=ticket-runner.target

[Service]
Type=simple
WorkingDirectory=%h/Documents/Programming/AI/ticket-runner
ExecStart=%h/.local/bin/node runner.js dashboard
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
# The dashboard is read-only and needs no coding-CLI auth or sandbox opt-out.
# It has no authentication, so 0.0.0.0 exposes ticket contents and repo paths to
# anyone who can reach this host on DASHBOARD_PORT. Set DASHBOARD_HOST=127.0.0.1
# to restrict it to loopback.
Environment=DASHBOARD_HOST=0.0.0.0
Environment=DASHBOARD_PORT=4600
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=ticket-runner.target
