# ticket-runner.service
[Unit]
Description=Notion AI ticket runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Documents/Programming/AI/ticket-runner
ExecStart=/usr/bin/node runner.js
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
